# ec2-ubuntu/setup/systemd-timers/yum-autoremove.service
[Unit]
Description=Run yum autoremove

[Service]
Type=oneshot
ExecStart=/usr/bin/yum autoremove --quiet
