# ec2/setup/systemd-timers/journalctl-vacuum.service
[Unit]
Description=Run journalctl vacuum

[Service]
Type=oneshot
ExecStart=/usr/bin/journalctl --vacuum-time=10d
